# lab/gateway-profiles/udm-pro-se-vm/initramfs/rootfs-payload/etc/systemd/system/unifi-stubd-vm-debug-ssh.service
[Unit]
Description=Install optional lab SSH debug public key
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/local/lib/unifi-stubd-vm/install-debug-ssh.sh
StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=multi-user.target
