# lab/gateway-profiles/udm-pro-se-vm/initramfs/rootfs-payload/etc/systemd/system/unifi-stubd-vm-netdevs.service
[Unit]
Description=Prepare QEMU lab network devices for UDM services
DefaultDependencies=no
After=systemd-udevd.service systemd-udev-trigger.service
Before=network-init.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/lib/unifi-stubd-vm/prepare-netdevs.sh
