# lab/gateway-profiles/udm-pro-se-vm/initramfs/rootfs-payload/etc/systemd/system/unifi-core.service.d/20-unifi-stubd-vm-simple.conf
[Service]
Type=simple
WatchdogSec=0
TimeoutStartSec=0
NotifyAccess=none
Environment=UNIFI_STUBD_VM_SYSTEM_INTERNET_READY=1
StandardOutput=journal+console
StandardError=journal+console
