# packaging/linux/usr/lib/systemd/system/unifi-stubd.service
[Unit]
Description=UniFi stub daemon
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/unifi-stubd
Restart=on-failure
RestartSec=5s
StateDirectory=unifi-stubd
ConfigurationDirectory=unifi-stubd

# Packaged defaults keep the advanced-adoption SSH shim closed. The service
# still stays non-root and receives only the low-port bind capability needed
# when an isolated lab explicitly enables ssh_listen on a privileged port.
User=unifi-stubd
Group=unifi-stubd
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/var/lib/unifi-stubd /etc/unifi-stubd

[Install]
WantedBy=multi-user.target
