# lab/gateway-profiles/udm-pro-se-vm/initramfs/rootfs-payload/etc/systemd/system/network-init.service.d/20-unifi-stubd-vm-netdevs.conf
[Unit]
Requires=unifi-stubd-vm-netdevs.service
After=unifi-stubd-vm-netdevs.service

[Service]
StandardOutput=journal+console
StandardError=journal+console
